# etc/crontab
0 * * * 3	/usr/local/bin/node /app/bin/lines.js update
*/1 * * * *	/usr/local/bin/node /app/bin/scores.js update
